# lib/misc/crontab
20 * * * * cd ~/current; bundle exec rake data:import_pdl

