# lib/misc/crontab
20 * * * * bundle exec rake data:import_pdl

